# system_files/shared/usr/lib/systemd/system/decky-install.service
[Unit]
Description=Install or update Decky Loader on first boot
Wants=network-online.target
After=network-online.target
ConditionFirstBoot=yes

[Service]
Type=oneshot
Environment=DECKY_INSTALL_MODE=stable
Environment=DECKY_INSTALL_MARKER_PATH=/var/lib/decky-installed
Environment=DECKY_ALLOW_ROOT=1
Environment=DECKY_INSTALL_USER=deck
ExecStart=/usr/libexec/bazzite-decky-install.sh

[Install]
WantedBy=multi-user.target
